# crontab
################ WIP - Need to create Docker cron job for currencies

1 * * * * ./cryptobot/main.py
# Mandatory blank line