# crontab
1 * * * * python ./cryptobot/main.py
# Mandatory blank line